# config/root_files/basilicata/etc/uci-defaults/99-hostname
#!/bin/sh

uci set system.@system[0].hostname='Openwrt-Ninux'
uci commit system

exit 0
